# wallet-daemon/deploy/walletd.service
# walletd — the 24/7 fedimint wallet daemon (spec §6a.8). Install as a systemd USER unit:
#   cargo install --path wallet-daemon  # installs `walletd` into ~/.cargo/bin
#   mkdir -p ~/.config/systemd/user
#   cp wallet-daemon/deploy/walletd.service ~/.config/systemd/user/
#   walletd init                       # scaffold config + 0600 token + client pointer
#   systemctl --user daemon-reload
#   systemctl --user enable --now walletd
#
# Fatal errors exit non-zero and systemd restarts the process; restart reconciles in-flight
# operations (abandon-and-resume is the model). SIGTERM triggers the ordered shutdown
# (stop intake → abort drivers → drain the actor).
[Unit]
Description=walletd — 24/7 fedimint wallet daemon + local API
Documentation=file://%h/.config/walletd/walletd.toml
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# `cargo install --path wallet-daemon` lands `walletd` in ~/.cargo/bin, which is NOT on the
# systemd user-manager's default PATH — a bare `ExecStart=walletd` would fail 203/EXEC. Use an
# absolute path (%h expands to the user's home); edit this if you installed the binary elsewhere.
ExecStart=%h/.cargo/bin/walletd
Restart=on-failure
RestartSec=2
# fedimint tunes RocksDB to a 2MB write buffer with NO extra memtable history, and its
# optimistic-commit validation needs history covering back to a transaction's snapshot. Its
# lnv2 receive task holds a transaction open across minutes-long long-polls, so ~2MB of
# client-side write churn in that window makes the commit die TryAgain -> WriteConflict ->
# panic, killing the shared receive task (the 24h-soak wedge; the settlement-stall watchdog
# restarts us when it happens). 32MB extends the survivable idle window 16x for a few tens
# of MB of daemon memory — the right trade for a 24/7 process.
Environment=FM_ROCKSDB_WRITE_BUFFER_SIZE=33554432
# The ordered shutdown is bounded (parked long-polls drain within their request deadline);
# give it generous headroom before systemd escalates to SIGKILL.
TimeoutStopSec=90

[Install]
WantedBy=default.target
